# deploy/d-brain-weekly.service
[Unit]
Description=d-brain Weekly Digest
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/life-pilot
ExecStart=/home/ubuntu/.local/bin/uv run python scripts/weekly.py
EnvironmentFile=/etc/systemd/system/d-brain.env
StandardOutput=journal
StandardError=journal
