# deploy/d-brain-bot.service
[Unit]
Description=D-Brain Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/life-pilot
ExecStart=/home/ubuntu/life-pilot/scripts/start-bot.sh
Restart=on-failure
RestartSec=10
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin:/bin
EnvironmentFile=/etc/systemd/system/d-brain.env

[Install]
WantedBy=multi-user.target
